# docker/php/supervisor/supervisord.conf
[supervisord]
logfile=/var/www/html/storage/logs/supervisord.log
logfile_maxbytes=5MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
user=root
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[program:horizon]
user=root
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan horizon
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/supervisor.log
stopwaitsecs=3600
stdout_logfile_maxbytes=5MB

[program:sentry-worker]
user=root
process_name=%(program_name)s_%(process_num)02d
command=docker-compose -f /home/docker-compose.yml run sentry sentry run worker

[program:sentry-cron]
user=root
process_name=%(program_name)s_%(process_num)02d
command=docker-compose -f /home/docker-compose.yml run sentry sentry run cron

[program:php-fpm]
user=root
command=php-fpm
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
